# root/etc/services.d/music_library_tools/run
#!/usr/bin/with-contenv sh

cd /music_library_tools
su abc -c "poetry run start"
