# root/etc/services.d/music_library_tools/run
#!/usr/bin/with-contenv sh

UMASK_SET=${UMASK_SET:-022}
umask "$UMASK_SET"

cd /music_library_tools

exec s6-setuidgid abc poetry run start
