# root/etc/services.d/music_library_tools/run
#!/usr/bin/with-contenv sh

cd /music_library_tools
poetry run start
