# root/etc/services.d/music_library_tools/run
#!/usr/bin/with-contenv sh

UMASK_SET=${UMASK_SET:-022}
umask "$UMASK_SET"

exec s6-setuidgid abc python3 /music_library_tools/music_library_tools/daemons.py
